# S3imggallery/puppet/files/s3imggal.conf
exec `which ruby2.0` /opt/company/utils/s3imggal/s3img.rb
start on startup
